# bastion-vm-firstboot/SOURCES/bastion-vm-firstboot.service
[Unit]
Description=Bastion VM first-boot setup (Homebrew + dev tools)
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/bastion-vm-firstboot/done

[Service]
Type=oneshot
TimeoutStartSec=infinity
User=user
ExecStartPre=+/usr/bin/install -d -o user -g user -m 0755 /var/lib/bastion-vm-firstboot
ExecStart=/usr/libexec/bastion-vm-firstboot/firstboot.sh
ExecStartPost=/usr/bin/touch /var/lib/bastion-vm-firstboot/done
StandardOutput=journal+console
StandardError=journal+console
RemainAfterExit=yes

# Hardening — selected to not break sudo (needed for corepack/dnf clean) and
# Homebrew's Ruby. See systemd-analyze security in CI for ongoing audit.
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
SystemCallArchitectures=native
# NOT set: NoNewPrivileges (sudo), ProtectSystem (dnf writes /var/cache),
# MemoryDenyWriteExecute (Homebrew Ruby), ProtectHome (we write ~/.npm-global,
# ~/.claude, /home/linuxbrew).

[Install]
WantedBy=multi-user.target
